# Autostart/blink-kiosk.service
[Unit]
Description=Blink Kiosk (Webcam + Web Server + Full Screen Browser)
After=network-online.target
Wants=network-online.target
# Ensure the desktop is ready before launching browser
After=graphical.target
Wants=graphical.target

[Service]
Type=simple
User=beta-blink
WorkingDirectory=/home/beta-blink/Blink_WebCam
Environment=HOME=/home/beta-blink

# Command sequence: start webcam, start Flask, then launch browser
ExecStart=/bin/bash -c '\
# Start webcam in background \
python3 Blink_WebCam_Capture.py & \
sleep 5; \
# Start Flask web server \
/home/beta-blink/.beta/bin/python3 Blink_Web_Server.py & \
sleep 5; \
# Launch browser full screen \
/usr/bin/google-chrome --kiosk http://wyse-beta-blink:5000 --noerrdialogs --disable-infobars --disable-session-crashed-bubble \
'

Restart=always

[Install]
WantedBy=graphical.target
